# builds/backend/crontab
*/5 * * * * root echo 'Cron command executed' >> /var/log/backend-cron.log
